# shell/supervisor/etc/supervisord.d/laravel-schedule.ini
[program:laravel-schedule]
process_name=%(program_name)s_%(process_num)02d
command=php /根目录/artisan schedule:run >> /dev/null 2>&1
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/项目根目录/storage/logs/schedule.log